# scripts/launchd/dev.jdagent.frontends.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.jdagent.frontends</string>

    <key>Comment</key>
    <string>JD Agent - Frontend Apps (Command Center, Tasks, Vault)</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jddavenport/Projects/JD Agent/scripts/start-frontends.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/jddavenport/Projects/JD Agent</string>

    <key>StandardOutPath</key>
    <string>/Users/jddavenport/Library/Logs/jd-agent-frontends.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jddavenport/Library/Logs/jd-agent-frontends.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jddavenport/.bun/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/jddavenport</string>
    </dict>

    <!-- Start after Docker is ready (60 second delay) -->
    <key>StartInterval</key>
    <integer>0</integer>

    <!-- Throttle restarts -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
