# scripts/launchd/dev.jdagent.docker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label: Unique identifier for this service -->
    <key>Label</key>
    <string>dev.jdagent.docker</string>

    <!-- Description -->
    <key>Comment</key>
    <string>JD Agent - Local Backend Services (Docker Compose)</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/docker</string>
        <string>compose</string>
        <string>-f</string>
        <string>/Users/jddavenport/Projects/JD Agent/docker-compose.yml</string>
        <string>up</string>
        <string>-d</string>
    </array>

    <!-- Run at system startup -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/jddavenport/Projects/JD Agent</string>

    <!-- Keep it running (restart if it crashes) -->
    <key>KeepAlive</key>
    <dict>
        <!-- Only keep alive if Docker is running -->
        <key>OtherJobEnabled</key>
        <dict>
            <key>com.docker.docker</key>
            <true/>
        </dict>
    </dict>

    <!-- Wait for network to be available -->
    <key>WaitForDebugger</key>
    <false/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/jddavenport/Library/Logs/jd-agent-docker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jddavenport/Library/Logs/jd-agent-docker.error.log</string>

    <!-- Environment variables (Docker needs PATH) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Start delay (wait for Docker Desktop to start) -->
    <key>StartInterval</key>
    <integer>0</integer>

    <!-- Throttle restarts (minimum seconds between launches) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
